# cdt-guard.timer
[Unit]
Description=Run Aliyun CDT Guard every minute

[Timer]
OnBootSec=30s
OnUnitActiveSec=60s
AccuracySec=5s
Persistent=true
Unit=cdt-guard.service

[Install]
WantedBy=timers.target
